# systemd/dvdarchiver-scan-consumer.service
# /etc/systemd/system/dvdarchiver-scan-consumer.service
[Unit]
Description=DVD Archiver - Consumer file de scan IA
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/dvdarchiver.conf
ExecStart=/usr/local/bin/scan_consumer.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
