# dvd-archiver/systemd/dvdarchiver-scan-consumer.path
[Unit]
Description=DVD Archiver - Surveillance Phase 2 (déclencheur scan)

[Path]
EnvironmentFile=/etc/dvdarchiver.conf
PathExistsGlob=${DEST:-/mnt/media_master}/*/raw/VIDEO_TS_BACKUP/VIDEO_TS/VIDEO_TS.VOB
Unit=dvdarchiver-scan-consumer.service

[Install]
WantedBy=multi-user.target
