# etc/systemd/system/dvd-ripd.service
#/etc/systemd/system/dvd-ripd.service 
[Unit]
Description=DVD rip queue daemon
After=network.target

[Service]
User=media
ExecStart=/usr/local/bin/dvd-ripd.sh
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
